# server/supervisord.conf
[unix_http_server]
file=/home/webdev/django/andrewramos/server/supervisord.sock

[supervisord]
logfile=/home/webdev/django/andrewramos/server/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/home/webdev/django/andrewramos/server/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/webdev/django/andrewramos/server/supervisord.sock

[program:andrewramos]
command=/home/webdev/.virtualenvs/andrewramos/bin/gunicorn_django -c server/gunicorn.conf.py
directory=/home/webdev/django/andrewramos
user=webdev
autostart=true
autorestart=true
stdout_logfile=/home/webdev/django/andrewramos/server/logs/supervisord.log
redirect_stderr=true